# dev/supervisord.conf
[program:miracle_server]
command = gunicorn miracle_server.wsgi --bind run/gunicorn.sock
stdout_logfile = logs/gunicorn_supervisor.log
redirect_stderr = true  